# gluu-flask.service
[Unit]
Description=Gluu FLask Daemon

[Service]
Type=forking
ExecStart=/usr/bin/gluuapi

[Install]
WantedBy=multi-user.target
