# gluuengine/templates/oxidp/memcached.conf
[program:memcached]
command=/usr/bin/memcached -p 11211 -u memcache -m 64 -t 4 -l 0.0.0.0 -vv
stdout_logfile=/var/log/memcached.log
stderr_logfile=/var/log/memcached.log
