# gluuengine/templates/ldap/opendj.conf
[program:opendj]
command=/opt/opendj/bin/start-ds --quiet -N
