# debian/postinst
#!/bin/sh -e
update-rc.d gluu-flask defaults > /dev/null
service gluu-flask start
